# files/eirka-conf/eirka-post.service
[Unit]
Description=Eirka Post
Requires=nginx.service mysql.service redis.service
After=nginx.service mysql.service redis.service

[Service]
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full

User=www-data
RuntimeDirectory=eirka
Environment=GIN_MODE=release
ExecStart=/data/pram/bin/eirka-post
Restart=on-failure

[Install]
WantedBy=multi-user.target
